# grpc/supervisord.conf
[supervisord]
nodaemon=true

[program:grpc_1]
command=/opt/grpc_1
autostart=true
autorestart=true

[program:grpc_2]
command=/opt/grpc_2
autostart=true
autorestart=true