# json_rpc/supervisord.conf
[supervisord]
nodaemon=true

[program:json_rpc_1]
command=/opt/json_rpc_1
autostart=true
autorestart=true

[program:json_rpc_2]
command=/opt/json_rpc_2
autostart=true
autorestart=true