# tcp/supervisord.conf
[supervisord]
nodaemon=true

[program:tcp_client]
command=/opt/tcp_client
autostart=true
autorestart=true

[program:tcp_server]
command=/opt/tcp_server
autostart=true
autorestart=true