# tcp/supervisord.conf
[supervisord]
nodaemon=true

[program:tcp_1]
command=/opt/tcp_1
autostart=true
autorestart=true

[program:tcp_2]
command=/opt/tcp_2
autostart=true
autorestart=true